# lecture-11-14/systemd/target/test.target
[Unit]
Description=Test Target
Conflicts=stop.target
