# lecture-11-12/systemd/intro/target/test.target
[Unit]
Description=Test Target

[Install]
WantedBy=graphical.target
